# playlist/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:fastapi]
command=/usr/bin/env uvicorn app:app --host 0.0.0.0 --port 8000
directory=/app   ; рабочая папка, где лежит app.py
user=root      ; под каким юзером запускать (могут быть root или ваш системный)
autostart=true ; автоматически стартовать при запуске supervisord
autorestart=true ; перезапускать, если процесс упадёт
startsecs=5    ; ждать 5 сек. после запуска, чтобы считать, что старт успешен
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:create_tables]
command=/bin/bash -c "sleep 30 && python create_tables.py"
directory=/app
user=root
autostart=true
autorestart=false
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
